# Server/private/crontab
* * * * * * * curl 127.0.0.1/cron >> /dev/null 2>&1