# Server/private/crontab
* * * * * * * curl 127.0.0.1/api/cron.php >> /var/log/cron.log 2>&1