# Server/private/crontab
* * * * * * * curl 127.0.0.1/api/cron.php >> /dev/null 2>&1